# debian/dump1090-fa/DEBIAN/postinst
#!/bin/bash
# postinst script for dump1090
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

NAME=dump1090-fa
RUNAS=dump1090
CRONFILE=/etc/cron.d/$NAME
TEMPLATECRON=/usr/share/$NAME/cron-template

case "$1" in
    configure)
        . /usr/share/debconf/confmodule

        if ! getent passwd "$RUNAS" >/dev/null
        then
            adduser --system --home /usr/share/$NAME --no-create-home --quiet "$RUNAS"
        fi

        # plugdev required for bladeRF USB access
        adduser "$RUNAS" plugdev

        # set up lighttpd
        if dpkg --compare-versions "$2" lt "3.1.0"
        then
            echo "Enabling lighttpd integration.." >&2
            lighty-enable-mod dump1090-fa || true

            # only enable the statcache config if there is nothing else around that already
            # configures it, because lighttpd fails if it's configured twice
            if ! grep -q -E '^\S*server.stat-cache-engine' /etc/lighttpd/conf-enabled/*.conf
            then
                echo "Enabling lighttpd integration (stat cache).." >&2
                lighty-enable-mod dump1090-fa-statcache || true
            fi
        fi

        # on upgrade, add an ENABLED line if it's not already present
        if dpkg --compare-versions "$2" lt-nl "3.7.0"
        then
            if [ -f /etc/default/dump1090-fa ]
            then
                if ! grep -q -E '^ENABLED=' /etc/default/dump1090-fa
                then
                    echo "Setting ENABLED=yes in /etc/default/dump1090-fa.." >&2
                    echo "# Automatically added by upgrade from $2" >>/etc/default/dump1090-fa
                    echo "ENABLED=yes" >>/etc/default/dump1090-fa
                fi
            fi
        fi

        if dpkg --compare-versions "$2" le "5.0"
        then
            echo "Enabling lighttpd skyaware module.." >&2
            lighty-enable-mod skyaware || true

        fi

        echo "Restarting lighttpd.." >&2
        invoke-rc.d lighttpd restart || true
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

# Automatically added by dh_systemd_enable/13.2~bpo10+1
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	# This will only remove masks created by d-s-h on package removal.
	deb-systemd-helper unmask 'dump1090-fa.service' >/dev/null || true

	# was-enabled defaults to true, so new installations run enable.
	if deb-systemd-helper --quiet was-enabled 'dump1090-fa.service'; then
		# Enables the unit on first installation, creates new
		# symlinks on upgrades if the unit file has changed.
		deb-systemd-helper enable 'dump1090-fa.service' >/dev/null || true
	else
		# Update the statefile to add new symlinks (if any), which need to be
		# cleaned up on purge. Also remove old symlinks.
		deb-systemd-helper update-state 'dump1090-fa.service' >/dev/null || true
	fi
fi
# End automatically added section
# Automatically added by dh_systemd_start/13.2~bpo10+1
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if [ -d /run/systemd/system ]; then
		systemctl --system daemon-reload >/dev/null || true
		deb-systemd-invoke start 'dump1090-fa.service' >/dev/null || true
	fi
fi
# End automatically added section
# Automatically added by dh_installinit/13.2~bpo10+1
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if [ -x "/etc/init.d/dump1090-fa" ]; then
		update-rc.d dump1090-fa defaults >/dev/null
		if [ -n "$2" ]; then
			_dh_action=restart
		else
			_dh_action=start
		fi
		invoke-rc.d dump1090-fa $_dh_action || exit 1
	fi
fi
# End automatically added section


if [ "$1" = "configure" ]; then db_stop; fi
exit 0
